# systemd/dotascraper.timer
[Unit]
Description=Run dotascraper every 30 minutes

[Timer]
OnCalendar=*-*-* *:00,30:00
RandomizedDelaySec=120

[Install]
WantedBy=timers.target